# radeon-profile-daemon/extra/radeon-profile-daemon.service
[Unit]
Description=radeon-profile daemon

[Service]
Type=simple
ExecStart=/usr/bin/radeon-profile-daemon

[Install]
WantedBy=multi-user.target
